# rc.local
/opt/bin/owampd -c /opt/etc-d /tmp -R /var/run -G nogroup -U nobody
